# scripts/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.quangdang.daily-digest-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/quangdang/Projects/AI-digest-v2/.venv/bin/python</string>
        <string>/Users/quangdang/Projects/AI-digest-v2/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/quangdang/Projects/AI-digest-v2</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/quangdang/Projects/AI-digest-v2/digest.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/quangdang/Projects/AI-digest-v2/digest_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/quangdang/Projects/AI-digest-v2/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/quangdang</string>
        <key>DIGEST_RUN_PROFILE</key>
        <string>publish</string>
        <key>USE_GROK_FOR_CLASSIFY</key>
        <string>1</string>
        <key>GROK_CLASSIFY_MODE</key>
        <string>retry</string>
        <key>USE_GROK_FOR_DELIVERY_RERANK</key>
        <string>1</string>
        <key>USE_GROK_FOR_FINAL_POLISH</key>
        <string>1</string>
        <key>MLX_MODEL</key>
        <string>mlx-community/Qwen2.5-14B-Instruct-4bit</string>
        <key>MLX_FALLBACK_MODEL</key>
        <string>mlx-community/Qwen2.5-14B-Instruct-4bit</string>
        <key>MLX_AUTO_FALLBACK</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
